# parity-rasp_1.3.0-0_armhf/DEBIAN/postinst
#!/bin/sh

systemctl enable parity
systemctl start parity

exit 0
